# deploy/systemd/cookster-queue.target
[Unit]
Description=Cookster Laravel Queue Workers
Wants=cookster-queue@1.service cookster-queue@2.service
